# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:orchestrator]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/orchestrator
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:email_agent]
command=uvicorn main:app --host 0.0.0.0 --port 8001
directory=/app/email_agent
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:calendar_agent]
command=uvicorn main:app --host 0.0.0.0 --port 8002
directory=/app/calendar_agent
autostart=true
autorestart=true
priority=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:task_agent]
command=uvicorn main:app --host 0.0.0.0 --port 8003
directory=/app/task_agent
autostart=true
autorestart=true
priority=4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:streamlit_ui]
command=streamlit run app.py --server.port=7860 --server.address=0.0.0.0
directory=/app/ui
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Made with Bob
